# config/systemd/user/pacsync.timer
[Unit]
Description=Runs pacsync once per day

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
