# config/systemd/user/pacsync.service
[Unit]
Description=Sync repositories and download package updates.

[Service]
Type=oneshot
ExecStart=/home/jake/.local/bin/pacsync
